# Assets/provisioners/hcl_domino_standalone_provisioner/0.1.23/scripts/provisioners/ansible_collections/startcloud/startcloud_roles/roles/web_terminal/templates/web-terminal.service.j2
[Unit]
Description=Web Terminal Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{ web_terminal_path }}
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=web-terminal

[Install]
WantedBy=multi-user.target
