# Assets/provisioners/hcl_domino_standalone_provisioner/0.1.23/scripts/provisioners/ansible_collections/startcloud/startcloud_roles/roles/bigbluebutton/templates/greenlight.service.j2
# /etc/systemd/system/greenlight.service
[Unit]
Description=Greenlight Service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
StandardError=null
StandardOutput=null
WorkingDirectory={{ bbb_greenlight_libdir }}
ExecStart=/usr/local/bin/docker-compose -f {{ bbb_greenlight_etcdir }}/docker-compose.yml up -d
ExecStop=/usr/local/bin/docker-compose -f {{ bbb_greenlight_etcdir }}/docker-compose.yml down

[Install]
WantedBy=multi-user.target
