# Assets/provisioners/hcl_domino_standalone_provisioner/0.1.23/scripts/provisioners/ansible_collections/startcloud/startcloud_roles/roles/rpi_lcd/templates/lcd-display.service.j2
[Unit]
Description=LCD Display
DefaultDependencies=no
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/display.py
StandardError=journal
StandardOutput=journal
Restart=on-abort

[Install]
WantedBy=sysinit.target