# Assets/provisioners/hcl_domino_additional_provisioner/0.1.23/provisioners/ansible_collections/startcloud/startcloud_roles/roles/wifi_connect/templates/wifi-connect.service.j2
[Unit]
Description=Wifi-Connect
DefaultDependencies=no
After=local-fs.target

[Service]
ExecStart=/usr/local/sbin/start-wifi-connect.sh
StandardError=journal
StandardOutput=journal

[Install]
WantedBy=sysinit.target