# Assets/provisioners/hcl_domino_standalone_provisioner/0.1.22/scripts/provisioners/ansible/roles/startcloud_tomcat/templates/tomcat.service.j2
[Unit]
Description=Tomcat {{ tomcat_version }}
After=syslog.target network.target

[Service]
Type=forking

User={{ tomcat_user }}
Group={{ tomcat_group }}


EnvironmentFile={{ tomcat_home_dir }}/environment
ExecStart=/usr/share/tomcat{{ tomcat_version }}/bin/catalina.sh start
ExecStop=/usr/share/tomcat{{ tomcat_version }}/bin/catalina.sh stop

[Install]
WantedBy=multi-user.target
