# Assets/provisioners/hcl_domino_additional_provisioner/0.1.23/scripts/provisioners/ansible_collections/startcloud/startcloud_roles/roles/rpi_nfc/templates/nfc-daemon.service.j2
[Unit]
Description=NFC daemon
After=network.target auditd.service

[Service]
User={{ service_user }}
ExecStart=/usr/bin/sudo /usr/local/bin/nfc-daemon
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartPreventExitStatus=255
Type=simple

[Install]
WantedBy=multi-user.target